# units/chat-api.service.d/credentials.conf
[Service]
LoadCredentialEncrypted=CHAT_DATABASE_URL:/opt/app/vault/credstore/CHAT_DATABASE_URL.cred
Environment=CHAT_DATABASE_URL_FILE=%d/CHAT_DATABASE_URL
LoadCredentialEncrypted=ACTIVITY_DATABASE_URL:/opt/app/vault/credstore/ACTIVITY_DATABASE_URL.cred
Environment=ACTIVITY_DATABASE_URL_FILE=%d/ACTIVITY_DATABASE_URL
LoadCredentialEncrypted=INTERNAL_SERVICE_JWT_SECRET:/opt/app/vault/credstore/INTERNAL_SERVICE_JWT_SECRET.cred
Environment=INTERNAL_SERVICE_JWT_SECRET_FILE=%d/INTERNAL_SERVICE_JWT_SECRET
